# etc/telegram.conf
description "Telegram-bot upstart script"

respawn
respawn limit 15 5

start on runlevel [2345]
stop on shutdown

setuid bot_dc
exec /bin/sh /Users/dray/git/telegram-bot/launch.sh
